# cookbooks/app_orient-latin.com/templates/default/init_solr.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          solr_<%= @app_url %>
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the <%= @app_url %> solr indexer
# Description: starts the <%= @app_url %> solr indexer
### END INIT INFO

directory="<%= @app_dir %>/current"
environment="production"
export HOME="/root"

source /usr/local/rvm/scripts/rvm
cd $directory >/dev/null
rvm <%= @app_rvm %>

start() {
  if ! running >/dev/null; then
    echo "Starting solr indexer..."
    RAILS_ENV=$environment rake sunspot:solr:start
  else
    echo "Solr is already running... pid=$(pid)"
  fi
}

stop() {
  if running >/dev/null; then
    echo "Stopping application..."
    RAILS_ENV=$environment rake sunspot:solr:stop
  else
    echo "Unable to stop, solr is not running..."
  fi
}

reindex() {
  if running >/dev/null; then
    RAILS_ENV=$environment rake sunspot:reindex
  else
    echo "Solr is not running..."
  fi
}

status() {
  if running; then
    echo "Solr is running... pid=$(pid)"
    exit 0
  else
    echo "Solr is stopped."
    exit 1
  fi
  return 
}

running() {
  pid >/dev/null 2>/dev/null
}

pid() {
  pgrep -f solr.data.dir=<%= @app_dir %>/releases/[0-9a-f]*/solr/data/$environment 2>/dev/null
}

usage() {
  echo "Usage: $0 start|stop|restart|reindex|status"
}

case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
restart)
  stop
  start
  ;;
reindex)
  reindex
  ;;
status)
  status
  ;;
*)
  usage >&2
  exit 1
  ;;
esac
