# cookbooks/graylog2/templates/default/graylog2.cron.erb
# This task is deployed by Chef !
*/10 * * * * root <%=@dir%>/rake_tasks.sh
