# cookbooks/graylog2/templates/default/graylog2.cron.erb
# This task is deployed by Chef !
*/10 * * * * root source /usr/local/rvm/scripts/rvm && cd <%=@dir%> && bundle exec rake streamalarms:send subscriptions:send
