# cookbooks/graylog2/templates/default/graylog2-web.init.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          <%= @app_url %>
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the <%= @app_url %> website/application
# Description: starts the <%= @app_url %> website/application
### END INIT INFO

directory="<%= @app_dir %>"
environment="production"
port="<%= @app_port %>"
export HOME="/root"

source /usr/local/rvm/scripts/rvm
cd $directory >/dev/null
rvm <%= @app_rvm %>

start() {
  if ! running >/dev/null; then
    echo "Starting application..."
    passenger start -a 127.0.0.1 -e $environment -p $port -d
  else
    echo "App is already running... pid=$(pid)"
  fi
}

stop() {
  if running >/dev/null; then
    echo "Stopping application..."
    passenger stop -p $port
  else
    echo "Unable to stop, app is not running..."
  fi
}

status() {
  if running; then
    echo "App is running... pid=$(pid)"
    exit 0
  else
    echo "App is stopped."
    exit 1
  fi
  return 
}

running() {
  ps -p $(pid) >/dev/null 2>/dev/null
}

pid() {
  cat $directory/tmp/pids/passenger.$port.pid 2>/dev/null
}

usage() {
  echo "Usage: $0 start|stop|restart|status"
}

case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
restart)
  stop
  start
  ;;
status)
  status
  ;;
*)
  usage >&2
  exit 1
  ;;
esac
